# shell/profile.d/05-path.sh
append_path "$HOME/.local/bin"
append_path "$HOME/go/bin"
