# shell/profile.d/10-editor.sh
# Set vim as the default editor
if command -v vim >/dev/null
then
	EDITOR=vim
fi
export EDITOR
