# shell/profile.d/20-ccache.sh
# gportay: https://wiki.archlinux.org/title/ccache#Enable_for_command_line

prepend_path "/usr/lib/ccache/bin"

# gportay: https://wiki.archlinux.org/title/ccache#Enable_with_colorgcc
# TODO append_path "/usr/lib/colorgcc/bin"
# TODO export CCACHE_PATH="/usr/bin"
