# shell/profile.d/30-buildroot.sh
# gportay: https://buildroot.org/downloads/manual/manual.html#env-vars
export BR2_DL_DIR="$HOME/.cache/buildroot/dl"
