# shell/profile.d/30-buildroot.sh
export BR2_DL_DIR="$HOME/.cache/buildroot/dl"
