# shell/profile.d/30-openembedded.sh
export DL_DIR="$HOME/.cache/oe/downloads"
export SSTATE_DIR="$HOME/.cache/oe/sstate-cache"
