# shell/profile.d/30-openembedded.sh
# gportay: https://docs.yoctoproject.org/ref-manual/variables.html#term-DL_DIR
export DL_DIR="$HOME/.cache/oe/downloads"
# gportay: https://docs.yoctoproject.org/ref-manual/variables.html#term-SSTATE_DIR
export SSTATE_DIR="$HOME/.cache/oe/sstate-cache"
