# shell/profile.d/85-nrf_connnect_sdk.sh
# gportay: https://developer.nordicsemi.com/nRF_Connect_SDK/doc/1.6.1/nrf/gs_installing.html#install-the-required-tools
#export PATH=${HOME}/gn:"$PATH"
