# shell/profile.d/95-ruby.sh
# gportay: https://wiki.archlinux.org/title/Ruby
GEM_HOME="$(gem env user_gemhome)"
export GEM_HOME
